# uwsgi_ng/files/etc_init_uwsgi.conf.jinja
# uwsgi

# AUTOGENERATED BY SALT. DO NOT EDIT

description "uwsgi service"
author "Marco Paolini <marco@elastic-coders.com>"

start on runlevel [2345]

stop on runlevel [016]

respawn

exec uwsgi --die-on-term {{ uwsgi_config }}

